# deploy/pm2-acoustic-monitor.service
[Unit]
Description=PM2 Acoustic Service Monitor
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/bin/bash /root/acoustic.uz/scripts/monitor-service-restart.sh
